# recipes-sample-ui/sample-ui-autorun/files/sample-ui.service
[Unit]
Description=Sample UI
After=mosquitto.service

[Service]
ExecStart=/usr/bin/run-sample-ui.sh

[Install]
WantedBy=multi-user.target
